# usr/lib/systemd/system/train-tickets-backend.service
[Unit]
Description=Train Tickets Backend API
After=systemd-tmpfiles-setup.service

[Service]
Type=simple
WorkingDirectory=/var/app/backend
ExecStart=/usr/bin/node /var/app/backend/src/index.js
Restart=on-failure
RestartSec=5

EnvironmentFile=-/etc/train-tickets/backend.env

[Install]
WantedBy=multi-user.target
