# usr/lib/systemd/system/train-tickets-backend.service
[Unit]
Description=Train Tickets Backend API
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/share/train-tickets/backend
ExecStart=/usr/bin/node /usr/share/train-tickets/backend/src/index.js
Restart=on-failure
RestartSec=5

EnvironmentFile=-/etc/train-tickets/backend.env

[Install]
WantedBy=multi-user.target
